# systemd/speakerctl.service
[Unit]
Description=Microsoft Modern USB-C Speaker button daemon
Documentation=https://github.com/interkelstar/microsoft-speaker-driver
After=systemd-udevd.service sound.target
Wants=sound.target

[Service]
Type=simple
User=speakerctl
Group=speakerctl
# belt-and-suspenders: groups also set via usermod in install.sh
SupplementaryGroups=input plugdev audio
ExecStart=/usr/local/bin/speakerctl
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=3s

# Security hardening
PrivateTmp=true
# NoNewPrivileges is intentionally NOT set: when [startup] pulse_user is
# configured, the daemon runs `sudo -u <user> pactl ...` to set PulseAudio
# volumes, which requires privilege escalation.
ProtectSystem=strict
ProtectHome=read-only
RuntimeDirectory=speakerctl

StandardOutput=journal
StandardError=journal
SyslogIdentifier=speakerctl

[Install]
WantedBy=multi-user.target
